# deployment/dash-ticket.service
[Unit]
Description=Dash Ticket Vendedor - Dashboard de Performance de Vendas
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=sup
WorkingDirectory=/home/sup/dash_ticket_vendedor
Environment=PATH=/home/sup/dash_ticket_vendedor/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
Environment=PYTHONPATH=/home/sup/dash_ticket_vendedor
EnvironmentFile=/home/sup/dash_ticket_vendedor/.env
ExecStart=/home/sup/dash_ticket_vendedor/venv/bin/gunicorn \
    --bind 0.0.0.0:8050 \
    --workers 2 \
    --timeout 120 \
    --access-logfile /home/sup/dash_ticket_vendedor/logs/access.log \
    --error-logfile /home/sup/dash_ticket_vendedor/logs/error.log \
    src.dashboard.app:server
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
