# deployment/dash-ticket-etl.service
[Unit]
Description=ETL Pipeline - Sankhya ERP para Data Warehouse (Diario D-1)
After=network.target dash-ticket.service

[Service]
Type=oneshot
User=sup
WorkingDirectory=/home/sup/dash_ticket_vendedor
Environment=PATH=/home/sup/dash_ticket_vendedor/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
Environment=PYTHONPATH=/home/sup/dash_ticket_vendedor
EnvironmentFile=/home/sup/dash_ticket_vendedor/.env
ExecStart=/home/sup/dash_ticket_vendedor/venv/bin/python /home/sup/dash_ticket_vendedor/src/data/etl/runner.py
StandardOutput=append:/home/sup/dash_ticket_vendedor/logs/etl.log
StandardError=append:/home/sup/dash_ticket_vendedor/logs/etl.log

[Install]
WantedBy=multi-user.target
