# deployment/dash-ticket-etl.timer
[Unit]
Description=ETL Pipeline Timer - Executa a cada 10 minutos
[Timer]
OnCalendar=*:0/10
Persistent=true

[Install]
WantedBy=timers.target
