# deployment/supervisord.conf
[supervisord]
nodaemon=false
logfile=/home/sup/dash_ticket_vendedor/logs/supervisord.log
pidfile=/home/sup/dash_ticket_vendedor/logs/supervisord.pid
childlogdir=/home/sup/dash_ticket_vendedor/logs

[unix_http_server]
file=/home/sup/dash_ticket_vendedor/logs/supervisor.sock

[supervisorctl]
serverurl=unix:///home/sup/dash_ticket_vendedor/logs/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ============================================================
; PROCESSO 1: Dash App (Gunicorn WSGI)
; ============================================================
[program:dash_ticket]
command=/home/sup/dash_ticket_vendedor/venv/bin/gunicorn
    --bind 0.0.0.0:8050
    --workers 2
    --timeout 120
    src.dashboard.app:server
directory=/home/sup/dash_ticket_vendedor
environment=PATH="/home/sup/dash_ticket_vendedor/venv/bin",PYTHONPATH="/home/sup/dash_ticket_vendedor"
autostart=true
autorestart=true
startretries=5
stdout_logfile=/home/sup/dash_ticket_vendedor/logs/gunicorn_out.log
stderr_logfile=/home/sup/dash_ticket_vendedor/logs/gunicorn_err.log

; ============================================================
; PROCESSO 2: ETL Runner (executa uma vez ao iniciar)
; Para agendamento diário, configure cron separado.
; ============================================================
[program:etl_startup]
command=/home/sup/dash_ticket_vendedor/venv/bin/python /home/sup/dash_ticket_vendedor/src/data/etl/runner.py
directory=/home/sup/dash_ticket_vendedor
environment=PATH="/home/sup/dash_ticket_vendedor/venv/bin",PYTHONPATH="/home/sup/dash_ticket_vendedor:/home/sup/dash_ticket_vendedor/src/data"
autostart=true
autorestart=false
startsecs=5
startretries=1
stdout_logfile=/home/sup/dash_ticket_vendedor/logs/etl_out.log
stderr_logfile=/home/sup/dash_ticket_vendedor/logs/etl_err.log
